# popcorn.cron
#!/bin/sh

if [ -x /usr/bin/popcorn-client -a -f /etc/popcorn.conf ]; then
    . /etc/popcorn.conf
    /usr/bin/popcorn-client
fi
