# popcorn.cron
#!/bin/sh

if [ -x /usr/bin/popcorn -a -r /etc/popcorn.conf ]; then
    . /etc/popcorn.conf
    export POPCORN_ENABLED
    export POPCORN_METHOD
    /usr/bin/popcorn
fi
